# target/linux/orion5x/base-files/etc/uci-defaults/99_remove_bridge

. /lib/functions.sh

board=$(board_name)

case "$board" in
dlink,dns323)
	# delete bridge device and set devices to 'eth0' directly
	uci delete network.@device[0]
	uci set network.lan.device='eth0'
	uci set network.lan6.device='eth0'
	uci commit network
	;;
esac
